# ansible/roles/step-ca-rootless-pkcs11/templates/systemd/p11-kit-server.socket.j2
[Unit]
Description=p11-kit server socket

[Socket]
ListenStream=%t/p11-kit/pkcs11-socket
SocketMode=0600

[Install]
WantedBy=sockets.target
